# player.service
[Unit]
Description=Wav player

[Service]
Type=simple
Exec=/usr/local/share/player/start.sh

[Install]
WantedBy=multi-user.target 